# webserver/supervisor.conf
[program:ia896]
command = /var/www/ia896/webserver.py
numprocs = 1
autorestart = true
user = root
stdout_logfile = /var/log/ia908.supervisor.log
redirect_stderr = true
